# slurmd.upstart
# slurmd - SLURM Resource Manager compute node daemon
#

description  "SLURM compute node daemon"

start on runlevel [23]
stop on runlevel [06S]

respawn

pre-start script
  MY_DAEMON="$(scontrol show daemons|sed -n 's/.*\(slurmd\).*/\1/p')"
  if [ "${MY_DAEMON}" != "slurmd" ] 
  then
    stop
    exit 0
  fi
end script

script
  SLURMD_OPTIONS=""
  /usr/sbin/slurmd -D ${SLURMD_OPTIONS}
end script

